# rtsp-cameras.service
[Unit]
Description=라즈베리파이 RTSP 웹캠 시스템
After=network.target
Wants=network.target

[Service]
Type=simple
User=tspol
Group=tspol
WorkingDirectory=/home/tspol/rsp-RTSP-cam
ExecStart=/usr/bin/python3 /home/tspol/rsp-RTSP-cam/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 환경 변수 설정
Environment=PYTHONPATH=/home/tspol/rsp-RTSP-cam
Environment=DISPLAY=:0

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/tspol/rsp-RTSP-cam

# 리소스 제한
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
